# docker/context/supervisord.conf
[supervisord]
nodaemon=true

[program:rtiexec]
command=/bin/bash ./start.sh
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[inet_http_server]
port=:9001
